# build/packaging/deb/DEBIAN/postinst
#!/bin/sh
set -e

echo "Running FIPS self-test..."
/usr/local/bin/cloudflared-selftest || {
    echo "WARNING: FIPS self-test failed. Check /usr/local/bin/cloudflared-selftest output."
    echo "The binary may not be FIPS-compliant on this system."
}

# Reload systemd if available
if command -v systemctl >/dev/null 2>&1; then
    systemctl daemon-reload || true
fi

echo ""
echo "cloudflared-fips installed successfully."
echo "  Config sample: /etc/cloudflared/config.yaml.sample"
echo "  Build manifest: /usr/share/cloudflared-fips/build-manifest.json"
echo "  Systemd unit: systemctl enable --now cloudflared"

#DEBHELPER#
